# config/usr/lib/systemd/system/zookeeper.service
[Unit]
Description=Starts Zookeeper service
After=network.target sshd.service

[Service]
EnvironmentFile=/etc/environment
Type=forking
ExecStart=/gdp/zookeeper/bin/zkServer.sh start
ExecStop=/gdp/zookeeper/bin/zkServer.sh stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
